# com.zzh.apple-music-weekly-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.zzh.apple-music-weekly-backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/zzh/Developer/音乐备份/backup_apple_music_weekly.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/zzh/Developer/音乐备份</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>2</integer>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>3</integer>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>20</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/zzh/Developer/音乐备份/launchd-apple-music-backup.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/zzh/Developer/音乐备份/launchd-apple-music-backup.error.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
